# systemd/agent-manager.service
[Unit]
Description=Agent Manager Daemon (HTTPS web password self-service + agent upgrade API)
Documentation=https://github.com/yjs-2026/agent-manager-daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# RuntimeDirectory= / StateDirectory= / LogsDirectory= / CacheDirectory=
# are systemd's first-class "create-if-missing" primitives. systemd
# creates them with the right ownership *before* ExecStart, which
# sidesteps the "Failed to set up mount namespacing" exit-226 issue
# we hit when using ReadWritePaths against a not-yet-created path.
# We override the defaults to match the names our config.yaml uses.
RuntimeDirectory=agent-manager
StateDirectory=agent-manager
LogsDirectory=agent-manager

# Listen on all interfaces (0.0.0.0) over HTTPS by default. Override
# AGENT_MANAGER_SERVER__BIND_HOST or AGENT_MANAGER_SERVER__BIND_PORT
# via an Environment= line below if you need to scope this down.
Environment="AGENT_MANAGER_CONFIG=/etc/agent-manager/config.yaml"
WorkingDirectory=/opt/agent-manager-daemon
ExecStart=/opt/agent-manager-daemon/.venv/bin/python -m agent_manager

# Hardening: we need /etc/shadow + chpasswd + systemctl, so a tight
# sandbox is impossible. Restrict what we *don't* need.
#
# MemoryDenyWriteExecute=false: CPython 3.12+ uses mmap(MAP_JIT) for
# certain code paths (notably cryptography import-time JIT, but also
# the in-process bytecode cache). MemoryDenyWriteExecute=true makes
# the daemon exit with code 226/NAMESPACE on first mmap. The cost of
# turning it off is small here because the daemon doesn't load
# untrusted code at runtime.
#
# No ReadWritePaths= needed: RuntimeDirectory=/StateDirectory=
# /LogsDirectory= above make those paths rw implicitly, and the
# install_root is mounted under /opt which ProtectSystem=strict leaves
# writable by default.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=false

Restart=on-failure
RestartSec=5s
# Cap the restart loop: if the daemon dies 5 times in 60s, give up
# and let `systemctl status` show the actual error instead of an
# endless "activating (auto-restart)" cycle.
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target